# user.crontab
# Backup Cloud Storage daily at midnight
0 0 * * * rsync -av --delete ~/cloud-storage/ ~/cloud-storage-daily-backup/

# Backup Cloud Storage weekly at midnight on Sunday
0 0 * * 0 rsync -av --delete ~/cloud-storage/ ~/cloud-storage-weekly-backup/
