# user.crontab
# Backup Cloud Storage weekly at midnight on Sunday
0 0 * * 0 ~/scripts/backup_cloud-storage
