# user.crontab
# Backup Cloud Storage daily at midnight
0 0 * * * ~/scripts/backup_cloud-storage
